# modules/mad/crontab
0 6 * * * root supervisorctl restart walker
