# modules/map-a-droid/crontab
0 6 * * * root supervisorctl restart walker
0 6 * * * root supervisorctl restart screenshot
